# ops/launchd/com.castlog.improvement-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.castlog.improvement-loop</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/user/ウェブサイト/lure-database/ops/run-agent.sh</string>
        <string>improvement-loop</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/user/ウェブサイト/lure-database</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/user/.nvm/versions/node/v22.22.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/user</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/user/ウェブサイト/lure-database/ops/logs/improvement-loop.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/ウェブサイト/lure-database/ops/logs/improvement-loop.error.log</string>
</dict>
</plist>
