# ops/launchd/com.castlog.oauth-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.castlog.oauth-check</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/user/ウェブサイト/lure-database/scripts/check-oauth-expiry.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/user/ウェブサイト/lure-database</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/user/ウェブサイト/lure-database/ops/logs/oauth-check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/ウェブサイト/lure-database/ops/logs/oauth-check.error.log</string>
</dict>
</plist>
