# ops/launchd/com.castlog.editorial-writer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.castlog.editorial-writer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/user/ウェブサイト/lure-database/ops/run-agent.sh</string>
        <string>editorial-writer</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/user/ウェブサイト/lure-database</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/user/.nvm/versions/node/v22.22.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/user/ウェブサイト/lure-database/ops/logs/editorial-writer.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/ウェブサイト/lure-database/ops/logs/editorial-writer.error.log</string>
</dict>
</plist>
